# main.service
[Unit]
Description = main micro-service

[Service]
ExecStart = /home/ubuntu/systemd/main